# systemd/dnsfwd.socket
[Unit]
Description=

[Socket]
ListenDatagram=53

[Install]
WantedBy=sockets.target
