# systemd/dnsfwd.socket
[Unit]
Description=

[Socket]
ListenDatagram=[::1]:53

[Install]
WantedBy=sockets.target
